# infrastructure/arkhe-os/99-arkhe-ignition.rules
# infrastructure/arkhe-os/99-arkhe-ignition.rules
# Trigger identity injection when USB with label ARKHE_KEYS is plugged
ACTION=="add", SUBSYSTEM=="block", ENV{ID_FS_LABEL}=="ARKHE_KEYS", RUN+="/usr/local/bin/arkhe-inject.sh %k"
